# supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/data/logs

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ============================================================
# PostgreSQL
# ============================================================
[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/postgres.log
stderr_logfile=/data/logs/postgres_error.log
startsecs=5

# ============================================================
# Redis
# ============================================================
[program:redis]
command=/usr/bin/redis-server --dir /data/redis --save 60 1 --appendonly yes --appendfilename "appendonly.aof"
autostart=true
autorestart=true
priority=20
stdout_logfile=/data/logs/redis.log
stderr_logfile=/data/logs/redis_error.log
startsecs=2

# ============================================================
# PostgreSQL backup (dumps to network volume every 10 minutes)
# ============================================================
[program:pg-backup]
command=/bin/bash -c "while true; do sleep 600; pg_dump -U postgres scalper > /data/postgres_backup/scalper.sql && echo \"$(date): backup ok\"; done"
user=postgres
autostart=true
autorestart=true
priority=25
stdout_logfile=/data/logs/pg_backup.log
stderr_logfile=/data/logs/pg_backup_error.log

# ============================================================
# SOL Scalper Trading Bot
# ============================================================
[program:scalper]
command=python3 -m src.main
directory=/root/sol-scalper
environment=HOME="/root",PATH="/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
priority=30
startsecs=10
startretries=5
stdout_logfile=/data/logs/scalper.log
stderr_logfile=/data/logs/scalper_error.log
# Wait for PostgreSQL and Redis to be ready
depends_on=postgres,redis
